# esp-mqtt-aqi/rootfs/etc/services.d/esp-mqtt-aqi/run
#!/usr/bin/with-contenv bashio

bashio::log.info 'Parameters:'
bashio::log.info 'Address:'
bashio::log.info "$(bashio::services mqtt 'host')"
bashio::log.info 'Port:'
bashio::log.info "$(bashio::services mqtt 'port')"
bashio::log.info 'Login:'
bashio::log.info "$(bashio::services mqtt 'username')"
bashio::log.info 'Password:'
bashio::log.info "$(bashio::services mqtt 'password')"

bashio::log.info 'Start esp-mqtt-aqi'

exec /usr/local/bin/esp-mqtt-aqi -d -h $(bashio::services mqtt "host") -u $(bashio::services mqtt "username") -p $(bashio::services mqtt "port") -P $(bashio::services mqtt "password")
