# upmpdcli/rootfs/etc/services.d/upmpdcli/run
#!/usr/bin/with-contenv bashio

bashio::log.info 'Start upmpdcli'
exec upmpdcli -c /etc/upmpdcli.conf