# upmpdcli/rootfs/etc/services.d/mpd/run
#!/usr/bin/with-contenv bashio

bashio::log.info 'Start MPD'
exec mpd --stderr --no-daemon /etc/mpd.conf